# infra/systemd/sc-idx-completeness-check.service
[Unit]
Description=SC_IDX price completeness check
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/Sustainacore
Environment=PYTHONPATH=/home/opc/Sustainacore
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/sustainacore/db.env
EnvironmentFile=/etc/sustainacore-ai/secrets.env
ExecStartPre=/usr/bin/test -r /etc/sustainacore/db.env
ExecStartPre=/usr/bin/test -r /etc/sustainacore-ai/secrets.env
ExecStartPre=/bin/ls -l /etc/sustainacore/db.env /etc/sustainacore-ai/secrets.env
ExecStartPre=/usr/bin/namei -l /etc/sustainacore/db.env
ExecStartPre=/usr/bin/namei -l /etc/sustainacore-ai/secrets.env
ExecStart=/usr/bin/flock -n /tmp/sc_idx_pipeline.lock /usr/bin/python3 /home/opc/Sustainacore/tools/index_engine/check_price_completeness.py --since-base --end today --email-on-fail
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
