# infra/systemd/sc-idx-pipeline.timer
[Unit]
Description=Run SC_IDX pipeline every 6 hours (UTC)

[Timer]
OnCalendar=*-*-* 00:30:00,05:30:00 UTC
Persistent=true
Unit=sc-idx-pipeline.service

[Install]
WantedBy=timers.target
