# infra/systemd/sc-idx-price-ingest.timer
[Unit]
Description=Run SC_IDX price ingest daily (23:30 UTC)

[Timer]
OnCalendar=*-*-* 23:30:00 UTC
Persistent=true
Unit=sc-idx-price-ingest.service

[Install]
WantedBy=timers.target
