# infra/systemd/sc-idx-price-ingest.timer
[Unit]
Description=Run SC_IDX price ingest every 6 hours (UTC)

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
OnCalendar=*-*-* 05:00:00 UTC
OnCalendar=*-*-* 09:00:00 UTC
OnCalendar=*-*-* 13:00:00 UTC
Persistent=true
Unit=sc-idx-price-ingest.service

[Install]
WantedBy=timers.target
