# infra/systemd/sc-idx-pipeline.service
[Unit]
Description=SC_IDX pipeline (ingest + completeness + impute + index calc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opc
WorkingDirectory=/opt/sustainacore-ai
Environment=PYTHONPATH=/opt/sustainacore-ai
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/sustainacore/db.env
EnvironmentFile=/etc/sustainacore-ai/secrets.env
ExecStartPre=/usr/bin/test -r /etc/sustainacore/db.env
ExecStartPre=/usr/bin/test -r /etc/sustainacore-ai/secrets.env
ExecStartPre=/bin/ls -l /etc/sustainacore/db.env /etc/sustainacore-ai/secrets.env
ExecStartPre=/usr/bin/namei -l /etc/sustainacore/db.env
ExecStartPre=/usr/bin/namei -l /etc/sustainacore-ai/secrets.env
ExecStart=/usr/bin/flock -n /tmp/sc_idx_pipeline.lock /usr/bin/python3 /opt/sustainacore-ai/tools/index_engine/run_pipeline.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
