# infra/systemd/sc-idx-completeness-check.timer
[Unit]
Description=Run SC_IDX completeness check weekdays (00:10 UTC)

[Timer]
OnCalendar=Mon..Fri *-*-* 00:10:00 UTC
Persistent=true
Unit=sc-idx-completeness-check.service

[Install]
WantedBy=timers.target
