# infra/systemd/sc-tech100-related.timer
[Unit]
Description=Tech100 related companies telemetry timer

[Timer]
OnCalendar=*-*-* 07:15:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
