# tools/research_generator/systemd/research-generator-manual.service
[Unit]
Description=Scheduled Research Generator (Manual Requests)

[Service]
Type=oneshot
User=opc
Group=opc
WorkingDirectory=/home/opc/Sustainacore
Environment=PYTHONPATH=/home/opc/Sustainacore:/home/opc/.local/lib/python3.9/site-packages:/home/opc/Sustainacore/tools/research_generator/.venv/lib/python3.9/site-packages
ExecStart=/usr/bin/python3 -m tools.research_generator.run_generator --process-manual --once

[Install]
WantedBy=multi-user.target
