# tools/research_generator/systemd/research-generator.service
[Unit]
Description=Scheduled Research Generator

[Service]
Type=oneshot
WorkingDirectory=/home/opc/Sustainacore
ExecStart=/home/opc/Sustainacore/tools/research_generator/.venv/bin/python -m tools.research_generator.run_generator --once

[Install]
WantedBy=multi-user.target
