# tools/research_generator/systemd/research-generator.timer
[Unit]
Description=Scheduled Research Generator Timer

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true
Unit=research-generator.service

[Install]
WantedBy=timers.target
