# tools/gemini_jobs_worker/systemd/gemini-jobs-worker.timer
[Unit]
Description=Gemini Jobs Worker Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=3min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
